# packaging/systemd/humen-mcp-self-update.service
[Unit]
Description=Update Humen MCP from the configured Arch package
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/humen-mcp-update.env
ExecStart=/usr/lib/humen-mcp/humen-mcp-self-update
